# root/etc/supervisord.d/gitlab.ini
[program: gitlab]
command = /assets/wrapper
priority=1
numprocs = 1
autostart = true
autorestart = true

