# root/etc/supervisord.d/shell.ini
[program:shell]
command=/bin/sh /shell/fixgitlab.sh
priority=5

