# root/etc/supervisord.d/sonar.ini
[program:sonarqube]
command=/bin/bash -c "/opt/sonarqube/bin/linux-x86-64/sonar.sh console"
stdout_events_enabled=true
stderr_events_enabled=true
